# opsworks_deploy_python/templates/default/shutdown-hook.conf.erb
description "run at shutdown"

start on starting rc
task
exec /bin/bash /etc/<%= @application %>/archive_<%= @application %>_logs.sh log
exec /bin/bash /etc/<%= @application %>/archive_supervisor_logs.sh log
