# opsworks_deploy_python/templates/default/shutdown-hook.conf.erb
description "run at shutdown"

start on starting rc
task
exec /bin/bash /etc/<%= @application %>/archive_logs.sh log
